# playbooks/config/fdb-prometheus-exporter.service
[Unit]
Description=Export fdb metrics to prometheus

[Service]
Type=simple
ExecStart=/usr/local/bin/fdb-prometheus-exporter.py /var/lib/node_exporter/fdb.prom
KillSignal=SIGINT
Restart=always
RestartSec=5s